# addon/roles/teleport_service/templates/teleport.service.j2
[Unit]
Description=Teleport Service

[Service]
Type=simple
ExecStart=teleport run --config=/root/.config/node_config.yaml
Restart=on-failure

[install]
WantedBy=multi-user.target